# dunst/.config/systemd/user/dunst.service
[Unit]
Description=Notification daemon

[Service]
ExecStart=/usr/bin/dunst
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

